# com.jasonos.observation-engine.dyson-hope.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.jasonos.observation-engine.dyson-hope</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/zsh</string>
        <string>/Volumes/SandboxData/code/observation-engine/run.sh</string>
        <string>--config</string>
        <string>/Volumes/SandboxData/code/observation-engine/configs/dyson-hope.yaml</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>8</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Volumes/SandboxData/Logs/observation-engine-dyson-hope.log</string>

    <key>StandardErrorPath</key>
    <string>/Volumes/SandboxData/Logs/observation-engine-dyson-hope-error.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
